# news-web.service
[Unit]
Description=News Now Web Application
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/news-web
Environment="PATH=/root/news-web/.venv/bin"
Environment="PYTHONPATH=/root/news-web"
ExecStart=/root/news-web/.venv/bin/gunicorn -c web/gunicorn.conf.py web.app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 