# news-web.service
[Unit]
Description=News Now Web Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/news-now
Environment="PATH=/var/www/news-now/.venv/bin"
Environment="PYTHONPATH=/var/www/news-now"
ExecStart=/var/www/news-now/.venv/bin/gunicorn -c web/gunicorn.conf.py -w 4 web.app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 