# new-app/app-aiworld.service
[Unit]
Description=app.aiworld.com.br Next.js application
After=network.target

[Service]
Type=simple
User=www-data
Environment=NODE_ENV=production
Environment=PORT=3000
WorkingDirectory=/var/www/app.aiworld.com.br
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/app.aiworld.com.br/out.log
StandardError=append:/var/log/app.aiworld.com.br/error.log
MemoryAccounting=true
MemoryHigh=300M
MemoryMax=400M

[Install]
WantedBy=multi-user.target
